# enet/rootfs/etc/services.d/enet-mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start enet service
# ==============================================================================

declare gateway
declare host
declare password
declare port
declare username

bashio::log.info "enet/run..."

# Config
gateway=$(bashio::config 'gateway')

# Services
host=$(bashio::services "mqtt" "host")
port=$(bashio::services "mqtt" "port")
username=$(bashio::services "mqtt" "username")
password=$(bashio::services "mqtt" "password")

# Wait until mosqitto is up and running
bashio::log.info "Waiting for mqtt ${host}:${port}..."
bashio::net.wait_for "${port}" "${host}"

# bashio::log.info "enet/run 4..."
export ENET_MQTT_PASS="${pass}"

# Run daemon
bashio::log.info "Starting enet-mqtt..."
bashio::log.info "enet/run 5..."
exec /usr/bin/enet-mqtt \
  --gateway "${gateway}" \
  --mqtt-host "${host}" \
  --mqtt-port "${port}" \
  --mqtt-user "${username}"
