# enet/rootfs/etc/services.d/enet-mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start enet service
# ==============================================================================

bashio::log.info "enet/run..."

ENET_GATEWAY=$(bashio::config 'gateway')
bashio::log.info "enet/run 1..."
MQTT_PASSWORD="$(bashio::jq /data/auth.json '.enet_password')"
bashio::log.info "enet/run 2..."

# Wait until mosqitto is up and running
bashio::log.info "Starting enet-mqtt..."
bashio::log.info "enet/run 3..."
bashio::net.wait_for 1883

bashio::log.info "enet/run 4..."
export ENET_MQTT_PASS="${MQTT_PASSWORD}"

# Run daemon
bashio::log.info "enet/run 5..."
exec /usr/bin/enet-mqtt \
  --gateway "${ENET_GATEWAY}" \
  --mqtt-host 127.0.0.1 \
  --mqtt-port 1883 \
  --mqtt-user "enet"
